# redis-backup.service
[Unit]
Description=Redis Backup to B2 Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/sh-rdb-b2
Environment=NODE_ENV=production
Environment=HOME=/root
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/bin/bash /root/sh-rdb-b2/backup-redis.sh
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=redis-backup

[Install]
WantedBy=multi-user.target 